# crawl_ctf.service
[Unit]
Description=Crawling capture the flag interface.

[Service]
WorkingDirectory=/home/wil/Desktop/crawl_ctf
ExecStart=/home/wil/Desktop/crawl_ctf/env/bin/python app.py
RestartSec=3s
Restart=always

[Install]
WantedBy=multi-user.target